# atlas.service
[Unit]
Description=ATLAS v2.0 - Telegram Intelligence System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/telbot
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /root/telbot/atlas_agent.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas

[Install]
WantedBy=multi-user.target
